# application/models/service/VP2-IP/crontab
# http://www.thegeekstuff.com/2011/07/php-cron-job/

MAILTO=alban.lopez+cron@gmail.com

01 */2 * * * cd /var/www/WsWds/application/models/service/VP2-IP/; php5 -f Archive.cron.php > log.txt

